# launchd/com.subtitle-translator.libre.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.subtitle-translator.libre</string>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>WorkingDirectory</key><string>/Users/sandeepinjarapu/.subtitle-translator</string>
  <key>ProgramArguments</key>
  <array>
    <string>/Users/sandeepinjarapu/.subtitle-translator/venv/bin/python3</string>
    <string>/Users/sandeepinjarapu/.subtitle-translator/local_translate_server.py</string>
  </array>
  <key>EnvironmentVariables</key>
  <dict>
    <key>HF_HOME</key><string>/Users/sandeepinjarapu/.subtitle-translator/hf-cache</string>
    <key>HF_HUB_DISABLE_XET</key><string>1</string>
    <key>HF_HUB_ENABLE_HF_TRANSFER</key><string>0</string>
  </dict>
  <key>StandardOutPath</key><string>/tmp/subtitle-translator-libre.log</string>
  <key>StandardErrorPath</key><string>/tmp/subtitle-translator-libre.log</string>
</dict>
</plist>
